# etc_udev_rules.d/61-onacpower.rules
# Rule for when switching to AC Power
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/local/share/power-profiles/performance.sh"

